# lixian
#! /bin/sh

### BEGIN INIT INFO
# Provides: lixian
# Required-Start: $local_fs $network $remote_fs $syslog
# Required-Stop: $local_fs $network $remote_fs $syslog
# Should-Start: network-manager
# Should-Stop: network-manager
# Default-Start:  2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: lixian daemon
# Description: sinosky lixian server daemon
### END INIT INFO

python2=/usr/local/bin/python2
lixianApp=/home/www/root/lixian/main.py
config=/home/www/etc/lixian.conf
log=/home/www/log/lixian.app.log
pidFile=/var/run/lixian.pid
if [ -f "$pidFile" ]; then
    pid=$(cat $pidFile)
    checkProgress=$(kill -0 $pid 2>&1 > /dev/null && echo -n $?)
    if [ "$checkProgress" != "0" ]; then
        rm -f $pidFile
        unset pid
    fi
fi
script=/etc/init.d/lixian

lixian() {
    nohup $python2 $lixianApp --f="$config" >> "$log" 2>&1 & echo -n $! > $pidFile
}

doStart() {
    if [ "$pid" ]; then
        echo "lixian already running. pid = $pid"
    else
        lixian
        echo 'Done.'
    fi
}

doStop() {
    if [ "$pid" ]; then
        kill $pid
        rm -f $pidFile
        echo 'Done.'
    else
        echo 'lixian not running.'
    fi
}

doRestart() {
    if [ "$pid" ]; then
        echo 'Stopping lixian...'
        kill $pid
        rm -f $pidFile
        echo -e 'Done.\nStarting lixian...'
        lixian
        echo 'Done.'
    else
        echo -e 'lixian not running.\nStarting lixian...'
        lixian
        echo 'Done.'
    fi
}

checkStatus() {
    if [ "$pid" ]; then
        echo "lixian is running. pid = $pid"
    else
        echo 'lixian not running.'
    fi
}

case "$1" in
    start)
    echo 'Starting lixian...'
    doStart
    ;;
    stop)
    echo 'Stopping lixian...'
    doStop
    ;;
    restart)
    echo 'Restarting lixian: '
    doRestart
    ;;
    status)
    checkStatus
    ;;
    *)
    echo "Usage: $script {start|stop|restart|status}" >&2
    exit 3
    ;;
esac

exit 0
